# deploy/systemd/resume-enhancer-nlp.service
[Unit]
Description=Resume Enhancer NLP Service
After=network.target
Requires=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/resume-enhancer/nlp
ExecStart=/home/ec2-user/resume-enhancer/nlp/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --workers 2
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=resume-enhancer-nlp

# Environment variables
Environment=PYTHONPATH=/home/ec2-user/resume-enhancer/nlp
EnvironmentFile=/home/ec2-user/resume-enhancer/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/home/ec2-user/resume-enhancer/logs

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=60%

[Install]
WantedBy=multi-user.target
